# deploy/ollama@.service
[Unit]
Description=Ollama LLM server (port %i)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu

# %i is the port number — enables: systemctl start ollama@11434 ollama@11435
Environment=OLLAMA_HOST=0.0.0.0:%i
Environment=HOME=/home/ubuntu

# Allow each process to handle 2 concurrent requests internally.
# N instances × 2 parallel = N×2 total concurrent slots.
Environment=OLLAMA_NUM_PARALLEL=2

ExecStart=/usr/local/bin/ollama serve

Restart=on-failure
RestartSec=5s
TimeoutStartSec=120

StandardOutput=journal
StandardError=journal
SyslogIdentifier=ollama-%i

[Install]
WantedBy=multi-user.target
