# deploy/rag-civil-smoke.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.civilsmartdictionary.smoke</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>REPO_DIR/deploy/smoke_check.sh</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>    <integer>9</integer>
    <key>Minute</key>  <integer>0</integer>
  </dict>

  <key>StandardOutPath</key>   <string>HOME_DIR/Library/Logs/rag-civil-smoke.log</string>
  <key>StandardErrorPath</key> <string>HOME_DIR/Library/Logs/rag-civil-smoke-err.log</string>
</dict>
</plist>
